# daemon/com.vex.lazarus.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vex.lazarus.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/valx/cathedral/1. AVos/1. Digital/projects/PROJECT_GODHAND_LAZARUS/venv/bin/python</string>
        <string>/Users/valx/cathedral/1. AVos/1. Digital/projects/PROJECT_GODHAND_LAZARUS/daemon/lazarus_sync_daemon.py</string>
        <string>--daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/valx/cathedral/1. AVos/1. Digital/projects/PROJECT_GODHAND_LAZARUS/daemon/lazarus_daemon.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/valx/cathedral/1. AVos/1. Digital/projects/PROJECT_GODHAND_LAZARUS/daemon/lazarus_daemon.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/valx/cathedral/1. AVos/1. Digital/projects/PROJECT_GODHAND_LAZARUS</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
